# releases/youtube-fuse-v2.1.0/youtube-fuse.service.template
[Unit]
Description=YouTube FUSE Filesystem
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=%PROJECT_DIR%
ExecStart=%PROJECT_DIR%/venv/bin/python youtube_api_fuse.py %MOUNT_POINT%
ExecStop=/bin/fusermount -u %MOUNT_POINT%
Restart=always
RestartSec=10
Environment=HOME=%HOME_DIR%

[Install]
WantedBy=multi-user.target
